# a2dp_receiver/avrcp.service
[Unit]
BindsTo=bluetooth.service
After=bluetooth.service
BindsTo=dev-ttyO3.device
After=dev-ttyO3.device

[Service]
ExecStart=/usr/bin/python /root/avrcp.py /dev/ttyO3

[Install]
WantedBy=multi-user.target
